# sysroots/i586-poky-linux-uclibc/var/lib/opkg/info/modutils-initscripts.list
/etc/init.d/modutils.sh
/etc/modules-load.d/
/etc/modules-load.d/bsp.conf
