# sysroots/i586-poky-linux-uclibc/var/lib/opkg/info/quark-init.list
/etc/init.d/quark-init.sh
/etc/modules-load.quark/galileo_gen2.conf
/etc/modules-load.quark/
/etc/modules-load.quark/galileo.conf
